# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
user=root

[program:tx5dr-server]
command=node /app/packages/server/dist/scripts/server-launcher.js /app/packages/server/dist/index.js
directory=/app
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=45
stopasgroup=true
killasgroup=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=NODE_ENV=production,PORT=4000,TX5DR_CONFIG_DIR=/app/data/config,TX5DR_DATA_DIR=/app/data,TX5DR_LOGS_DIR=/app/data/logs,TX5DR_CACHE_DIR=/app/data/cache,GLIBC_TUNABLES="glibc.rtld.execstack=2",RTC_DATA_AUDIO_UDP_PORT="%(ENV_RTC_DATA_AUDIO_UDP_PORT)s",RTC_DATA_AUDIO_ICE_UDP_MUX="%(ENV_RTC_DATA_AUDIO_ICE_UDP_MUX)s"
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
