# linux/tx5dr.service
[Unit]
Description=TX-5DR Digital Radio Server
After=network.target nginx.service
Wants=nginx.service

[Service]
Type=simple
User=tx5dr
Group=tx5dr
EnvironmentFile=-/etc/tx5dr/config.env
# Allow loading shared objects with executable stack (PT_GNU_STACK RWE).
# Required by wsjtx-lib's Fortran code (libwsjtx_core.so) for trampoline support.
# glibc >= 2.41 (Debian trixie) blocks this by default; value 2 = always allow.
Environment=GLIBC_TUNABLES=glibc.rtld.execstack=2
ExecStart=/usr/bin/node /usr/share/tx5dr/packages/server/dist/scripts/server-launcher.js /usr/share/tx5dr/packages/server/dist/index.js
WorkingDirectory=/usr/share/tx5dr
Restart=on-failure
RestartSec=5
TimeoutStopSec=45s
StandardOutput=journal
StandardError=journal

# Security hardening
# Note: NoNewPrivileges=true is intentionally omitted because wsjtx-lib's
# Fortran code (libwsjtx_core.so) requires executable stack (PT_GNU_STACK RWE)
# for trampoline support. NoNewPrivileges prevents the kernel from enabling
# executable stack, causing dlopen to fail with EINVAL.
ProtectSystem=strict
ReadWritePaths=/var/lib/tx5dr

[Install]
WantedBy=multi-user.target
